# etc/systemd/system/ntserv.service
[Unit]
Description=Sanic Example

[Service]
User=shane
ExecStart=/usr/bin/env python3 -m ntserv
Restart=always
Environment="ENV=prod" "WORKERS=2" "PSQL_PASSWORD=REDACTED" "JWT_SECRET=REDACTED" "PROXY_SECRET=REDACTED"

[Install]
WantedBy=multi-user.target

